# auto-update.timer
[Unit]
Description=auto-update system timer

[Timer]
After=network-online
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
