# ublue-updater.timer
[Unit]
Description=update-ublue system timer

[Timer]
After=network-online
OnUnitActiveSec=24h

[Install]
WantedBy=timers.target
