# files/usr/lib/systemd/user/ublue-update.service
[Unit]
Description=Universal Blue Update Oneshot Service

[Service]
Restart=on-failure
RestartSec=1h
Type=oneshot
ExecStart=/usr/bin/ublue-update


