# files/usr/lib/systemd/user/ublue-update.service

[Unit]
Description=auto-update oneshot service

[Service]
Type=oneshot
ExecStart=/usr/bin/ublue-update
Restart=on-failure
RestartSec=15min
