# files/usr/lib/systemd/system/ublue-update.service
[Unit]
Description=Universal Blue Update Oneshot Service

[Service]
Type=oneshot
ExecStart=/usr/bin/ublue-update
